# 83-usb2host.rules.in
ACTION=="remove", SUBSYSTEM=="typec", RUN+="%BINARY% %NODE%"