# usb2host.service.in
[Unit]
Description=Init dual-role USB port to host mode

[Service]
Type=oneshot
ExecStart=%BINARY% %NODE%

[Install]
WantedBy=multi-user.target